# distro/customer-vps/systemd/matrix-homeserver.service
[Unit]
Description=Matrix OS messaging homeserver (Synapse)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=matrix
Group=matrix
Environment=MATRIX_MESSAGING_ROOT=/mnt/HC_Volume_104683898/matrix-messaging
WorkingDirectory=/opt/matrix/messaging
StateDirectory=matrix-messaging
ExecStart=/opt/matrix/messaging/bin/synapse --config-path ${MATRIX_MESSAGING_ROOT}/synapse/homeserver.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
