# distro/customer-vps/systemd/matrix-db-backup.timer
[Unit]
Description=Run Matrix OS customer database backup hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target

