# distro/customer-vps/systemd/matrix-bridge-whatsapp.service
[Unit]
Description=Matrix OS WhatsApp bridge
After=network-online.target matrix-homeserver.service
Wants=network-online.target matrix-homeserver.service

[Service]
Type=simple
User=matrix
Group=matrix
Environment=MATRIX_MESSAGING_ROOT=/mnt/HC_Volume_104683898/matrix-messaging
WorkingDirectory=/opt/matrix/messaging
ExecStart=/opt/matrix/messaging/bin/mautrix-whatsapp -c ${MATRIX_MESSAGING_ROOT}/mautrix-whatsapp/config.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
